# installer/telegram-saas.service.example
[Unit]
Description=Telegram SaaS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/telegram-saas
ExecStart=/root/telegram-saas/venv/bin/python /root/telegram-saas/telegram_api_v4.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
